# piHeatingLCD/piHeatingLCD/gpio.service
[Unit]
 Description=Read GPIO Service
 After=multi-user.target

 [Service]
 User=pi
 Group=pi
 Type=idle
 #ExecStart=/usr/bin/gpio-watch -s /home/pi/piHeatingLCD/gpio-scripts -e switch 18 23 24 25
 ExecStart=/usr/bin/gpio-watch -s /home/pi/piHeatingLCD/gpio-scripts -e switch 18 23 24 25 -l /home/pi/piHeatingLCD/gpio-watch.log > /dev/null 2>&1

 [Install]
 WantedBy=multi-user.target
 
